# base-skel/etc/profile.d/epics-base.sh
export PATH="$PATH:/opt/epics/base/bin/linux-x86_64"
